# ostree-notify/ostree-notify.timer
[Unit]
Description=rpm-ostree notification of staged deployment

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
